# atomupd-daemon/atomupd.service.in
[Unit]
Description=Atomic Update Daemon

[Service]
Type=dbus
BusName=com.steampowered.Atomupd1
ExecStart=@libexecdir@/atomupd-daemon
